# systemd/clp-pm2-autoinstall.path
[Unit]
Description=Watch nginx sites-enabled for new CloudPanel sites (PM2 auto-install)

[Path]
PathChanged=/etc/nginx/sites-enabled
Unit=clp-pm2-autoinstall.service

[Install]
WantedBy=multi-user.target
